# deploy/supervisor/pi/checkit.conf
[program:checkit_backend]
directory=/home/pi/CheckIT/backend
command=/home/pi/CheckIT/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
autostart=true
autorestart=true
stderr_logfile=/home/pi/CheckIT/backend/backend_err.log
stdout_logfile=/home/pi/CheckIT/backend/backend_out.log
user=pi
environment=PATH="/home/pi/CheckIT/backend/venv/bin:%(ENV_PATH)s",PYTHONPATH="/home/pi/CheckIT/backend"

[program:checkit_frontend]
directory=/home/pi/CheckIT/frontend
command=npm run dev -- --host 0.0.0.0
autostart=true
autorestart=true
stderr_logfile=/home/pi/CheckIT/frontend/frontend_err.log
stdout_logfile=/home/pi/CheckIT/frontend/frontend_out.log
user=pi
environment=PATH="/usr/local/bin:/usr/bin:/bin:%(ENV_PATH)s"
