# deploy/supervisor/pi/checkit.conf
[program:checkit_backend]
directory=/root/CheckIT
command=bash /root/CheckIT/start_rpi.sh
autostart=true
autorestart=true
stderr_logfile=/root/CheckIT/backend_err.log
stdout_logfile=/root/CheckIT/backend_out.log
user=root
environment=PATH="/usr/local/bin:/usr/bin:/bin:%(ENV_PATH)s"

[program:checkit_frontend]
directory=/root/CheckIT/frontend
command=npm run dev -- --host 0.0.0.0
autostart=true
autorestart=true
stderr_logfile=/root/CheckIT/frontend/frontend_err.log
stdout_logfile=/root/CheckIT/frontend/frontend_out.log
user=root
environment=PATH="/usr/local/bin:/usr/bin:/bin:%(ENV_PATH)s"
